# deploy/configuration/crons/get_and_parse.cron
#!/bin/bash

set -e

/home/ubuntu/gitroot/pancancer-info/pancan-scripts/run_get.sh > /home/ubuntu/gitroot/pancancer-info/pancan-scripts/run_get.sh.log 2>&1
/home/ubuntu/gitroot/public-workflows/decider-bwa-pancancer/parse.sh > /home/ubuntu/gitroot/public-workflows/decider-bwa-pancancer/parse.sh.log 2>&1

